# home/core/.config/systemd/user/acme-sh-cron.timer
[Unit]
Description=Acme.sh cron timer

[Timer]
OnUnitActiveSec=1day
Unit=acme-sh-cron.service
