# home/core/.config/systemd/user/backup.timer
[Unit]
Description=Backup timer

[Timer]
OnCalendar=*-*-* 02:00:00
Unit=duplicity-backup.service

[Install]
WantedBy=timers.target
