# home/core/.config/systemd/user/cloudflare-ddns.timer
[Unit]
Description=Cloudflare DDNS timer

[Timer]
OnActiveSec=0
OnUnitActiveSec=5min
Unit=cloudflare-ddns.service
