# home/core/.config/systemd/user/nextcloud-cron.timer
[Unit]
Description=Nextcloud Cron timer

[Timer]
OnActiveSec=0
OnUnitActiveSec=5min
Unit=nextcloud-cron.service
