# home/core/.config/systemd/user/nextcloud-cron.service
[Unit]
Description=Nextcloud cron service

[Service]
Type=oneshot
ExecStart=podman exec -it nextcloud /lifecycle.sh periodic
